# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root
logfile_maxbytes=50MB
logfile_backups=10

; ======================================
; Couchbase 服务
; ======================================
[program:couchbase]
command=/opt/couchbase/bin/couchbase-server
directory=/opt/couchbase/var
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/couchbase-stderr.log
stdout_logfile=/var/log/supervisor/couchbase-stdout.log
priority=10
startsecs=60
stopwaitsecs=30

; ======================================
; Couchbase 初始化
; ======================================
[program:couchbase-init]
command=/opt/couchbase/init-couchbase.sh
user=root
autostart=true
autorestart=false
startretries=0
stderr_logfile=/var/log/supervisor/couchbase-init-stderr.log
stdout_logfile=/var/log/supervisor/couchbase-init-stdout.log
priority=11
startsecs=0

; ======================================
; Node.js 应用 (通过 PM2 管理)
; ======================================
[program:nodejs-apps]
command=pm2-runtime start /app/ecosystem.config.js --no-daemon
directory=/app
user=nodejs
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nodejs-apps-stderr.log
stdout_logfile=/var/log/supervisor/nodejs-apps-stdout.log
priority=20
startsecs=30
environment=NODE_ENV="production"
